# deployment/puppeteer.service
[Unit]
Description=Hive Mind Puppeteer Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/hive-mind
Environment=NODE_ENV=production
Environment=GOOGLE_CLIENT_ID=your_client_id
Environment=GOOGLE_CLIENT_SECRET=your_client_secret
Environment=GOOGLE_REDIRECT_URI=your_redirect_uri
Environment=GOOGLE_REFRESH_TOKEN=your_refresh_token
Environment=DATABASE_URL=postgresql://user:password@affine-vm:5432/hivemind
ExecStart=/usr/bin/node src/puppeteer-service.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 